# docker/blog.container
[Unit]
Description=Astro Blog Container
After=network.target

[Container]
Image=blog:latest
ContainerName=blog
PublishPort=4321:80
AutoUpdate=registry
HealthCmd=curl -f http://localhost/ || exit 1
HealthInterval=30s
HealthTimeout=3s
NoNewPrivileges=true

[Service]
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target default.target